# bin/rc.local
#!/bin/bash

chmod 777 /tmp
/home/krgr/bin/ramdisk.start &
